# spec/fixtures/dev.dnsmasq.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.dnsmasq</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>Redirects *.dev DNS requests to localhost.</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/test/boxen/homebrew/sbin/dnsmasq</string>
      <string>--conf-file=/test/boxen/config/dnsmasq/dnsmasq.conf</string>
      <string>--keep-in-foreground</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <dict>
      <key>NetworkState</key>
      <true/>
    </dict>

    <key>WorkingDirectory</key>
    <string>/test/boxen/data/dnsmasq</string>

    <key>StandardErrorPath</key>
    <string>/test/boxen/log/dnsmasq/console.log</string>

    <key>StandardOutPath</key>
    <string>/test/boxen/log/dnsmasq/console.log</string>
  </dict>
</plist>
