# templates/dev.dnsmasq.plist.erb
<?xml version="1.0" encoding="UTF-8"?>

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.dnsmasq</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>Redirects *.dev DNS requests to localhost.</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string><%= scope.lookupvar "dnsmasq::config::executable" %></string>
      <string>--conf-file=<%= scope.lookupvar "dnsmasq::config::configfile" %></string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>WorkingDirectory</key>
    <string><%= scope.lookupvar "dnsmasq::config::datadir" %></string>

    <key>StandardErrorPath</key>
    <string><%= scope.lookupvar "dnsmasq::config::logfile" %></string>

    <key>StandardOutPath</key>
    <string><%= scope.lookupvar "dnsmasq::config::logfile" %></string>
  </dict>
</plist>
